# infra/talkeo-api.service
[Unit]
Description=Talkeo API (stockage VPS SQLite)
After=network.target
ConditionPathExists=/opt/crypt/src/server/src/index.js

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/crypt/src/server
# Le tiret (-) évite le crash systemd si .env absent au premier boot
EnvironmentFile=-/opt/crypt/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=talkeo-api

[Install]
WantedBy=multi-user.target
