# infra/talkeo-api.service
[Unit]
Description=Talkeo API (stockage VPS SQLite)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/crypt/src/server
EnvironmentFile=/opt/crypt/.env
Environment=NODE_ENV=production
ExecStart=/bin/bash -c 'exec $(which node || echo /usr/bin/node) src/index.js'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=talkeo-api

[Install]
WantedBy=multi-user.target
